# docker-all-in-one/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
startsecs=0
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
numprocs=1
priority=10
autostart=true
autorestart=true
redirect_stderr=true

[program:mongod]
command=mongod
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:server]
command=node /data/server-app/server.js
autostart=true
autorestart=true
environment=NODE_ENV=production
user=app
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:audit]
command=node /data/audit-app/server.js
autostart=true
autorestart=true
environment=NODE_ENV=production
user=app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:servermigration]
command=node /data/server-app/node_modules/mongodb-migrate-maestroserver -runmm -cfg .migration-config.js up 
directory=/data/server-app
autostart=true
autorestart=false
environment=NODE_ENV=production
user=app
startretries=0
numprocs=1
startsecs=10
redirect_stderr=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:discovery]
command=/usr/local/bin/gunicorn --config /data/discovery-app/gunicorn_config.py --chdir /data/discovery-app/ run:app
directory=/data/discovery-app
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:discoveryworker]
command=bash -c 'sleep 10 && celery -A app.celery worker -l info --hostname=discovery@%%h'
directory=/data/discovery-app/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/discovery-data/
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:report]
command=/usr/local/bin/gunicorn --config /data/report-app/gunicorn_config.py --chdir /data/report-app/ run:app
directory=/data/report-app/
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:reportworker]
command=bash -c 'sleep 10 && celery -A app.celery worker -l info --hostname=report@%%h'
directory=/data/report-app/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/report-data/
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:analytics]
command=/usr/local/bin/gunicorn --config /data/analytics-maestro/gunicorn_config.py --chdir /data/analytics-maestro/ run:app
directory=/data/analytics-maestro/
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:analyticsworker]
command=bash -c 'sleep 10 && celery -A app.celery worker -l info --hostname=analytics@%%h'
directory=/data/analytics-maestro/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/analytics-maestro/
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:analyticsfront]
command=node /data/analytics-front/server.js
autostart=true
autorestart=true
environment=NODE_ENV=production
user=app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log

[program:data]
command=/usr/local/bin/gunicorn --config /data/data-app/gunicorn_config.py --chdir /data/data-app/ run:app
directory=/data/data-app
user=app
autostart=true
priority=5
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log