# docker-all-in-one/supervisordWithoutStore.conf
[supervisord]
nodaemon=true


[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=0
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
priority=5

[program:server]
command=node /data/server-app/server.js
directory=/data/server-app
autostart=true
autorestart=true
environment=NODE_ENV=production
user=app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=10

[program:audit]
command=node /data/audit-app/server.js
autostart=true
autorestart=true
environment=NODE_ENV=production
user=app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=50

[program:discovery]
command=/home/app/venv/bin/gunicorn --config /data/discovery-app/gunicorn_config.py --chdir /data/discovery-app/ run:app
directory=/data/discovery-app
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
environment=PATH="/home/app/venv/bin:%(ENV_PATH)s"
priority=50

[program:discoveryworker]
command=celery -A app.celery worker -l info --hostname=discovery@%%h
directory=/data/discovery-app/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/discovery-data/,PATH="/home/app/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=100

[program:report]
command=/home/app/venv/bin/gunicorn --config /data/report-app/gunicorn_config.py --chdir /data/report-app/ run:app
directory=/data/report-app/
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
environment=PATH="/home/app/venv/bin:%(ENV_PATH)s"
priority=50

[program:reportworker]
command=celery -A app.celery worker -l info --hostname=report@%%h
directory=/data/report-app/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/report-data/,PATH="/home/app/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=100

[program:analytics]
command=/home/app/venv/bin/gunicorn --config /data/analytics-maestro/gunicorn_config.py --chdir /data/analytics-maestro/ run:app
directory=/data/analytics-maestro/
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
environment=PATH="/home/app/venv/bin:%(ENV_PATH)s"
priority=50

[program:analyticsworker]
command=celery -A app.celery worker -l info --hostname=analytics@%%h
directory=/data/analytics-maestro/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/analytics-maestro/,PATH="/home/app/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=100

[program:analyticsfront]
command=node /data/analytics-front/server.js
autostart=true
autorestart=true
environment=NODE_ENV=production
user=app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=50

[program:scheduler]
command=python /data/scheduler-app/run.py
directory=/data/scheduler-app/
user=app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
environment=PATH="/home/app/venv/bin:%(ENV_PATH)s"
priority=50

[program:schedulerworker]
command=celery -A app.celery worker -l info  -S app.schedulers.MongoScheduler --hostname=scheduler@%%h
directory=/data/scheduler-app/
user=app
autostart=true
autorestart=true
environment=CELERYD_CHDIR=/data/scheduler-data/,PATH="/home/app/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=100

[program:data]
command=/home/app/venv/bin/gunicorn --config /data/data-app/gunicorn_config.py --chdir /data/data-app/ run:app
directory=/data/data-app
user=app
autostart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
environment=PATH="/home/app/venv/bin:%(ENV_PATH)s"
priority=10

[program:ws]
command=/usr/local/bin/centrifugo -c /config.json
user=app
autostart=true
priority=5
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
priority=30

[program:servermigration]
command=node /data/server-app/node_modules/mongodb-migrate-maestroserver -runmm -cfg .migration-config.js up
directory=/data/server-app
autostart=true
autorestart=false
environment=NODE_ENV=production
user=app
startretries=0
numprocs=1
startsecs=20
redirect_stderr=false
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
priority=999